# crontab
* 0 * * * root bash /strip/stripper.sh >> /var/log/cron.log
