# crontab
* 0 * * * root sh /strip/stripper.sh >> /var/log/cron.log
